# src/deb/systemd/rebazer.service
[Unit]
Description=ReTest ReBaZer

[Service]
ExecStart=/usr/bin/rebazer
Type=simple
User=rebazer

[Install]
WantedBy=multi-user.target
