# src/deb/control/postinst
#!/bin/sh

systemctl daemon-reload 
systemctl start [[name]].service 
systemctl enable [[name]].service
